# board/raspberrypi/overlay/etc/udev/rules.d/80-alsa.rules
KERNEL=="controlC[0-7]", ACTION=="add|change", RUN+="alsa-utils"
